# monero/cpu/etc/s6/dnsmasq/run
#!/bin/sh

#### exec dnsmasq -C /etc/dnsmasq.conf